# onedrive@.service.in
[Unit]
Description=OneDrive Free Client
Documentation=https://github.com/skilion/onedrive
After=network.target

[Service]
ExecStart=@PREFIX@/bin/onedrive -m --confdir="~/.config/%i"
Restart=no

[Install]
WantedBy=default.target
