# onedrive.service.in
[Unit]
Description=OneDrive Free Client
Documentation=https://github.com/skilion/onedrive

[Service]
ExecStart=@PREFIX@/bin/onedrive -m
Restart=on-failure

[Install]
WantedBy=default.target
