# services/upstart.conf
# onedrive.conf - OneDrive Upstart service
# Documentation at https://github.com/skilion/onedrive

description "OneDrive Free Client"

start on started network-services
stop on stopping network-services

exec /usr/local/bin/onedrive -m
